# service/docksec.service
[Unit]
Description=docksec
Before=docker.service
After=network.target

[Service]
ExecStart=/usr/bin/docksec

[Install]
WantedBy=multi-user.target
